# com.slap.runner.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.slap.runner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/slap/.venv/bin/python</string>
        <string>/ABSOLUTE/PATH/TO/slap/slap.py</string>
        <string>runner</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/slap</string>

    <!-- REFERENCE ONLY — the real plist is generated by `python slap.py plist`,
         which reads config.yaml's schedule.active_days and fire_window_start
         and emits one StartCalendarInterval entry per active day (shown here
         for Mon-Fri, matching config.yaml.example's default). A single
         StartCalendarInterval dict can't express "these weekdays only" —
         launchd's own answer is an ARRAY of dicts, one per day (Weekday
         1=Mon...5=Fri, 6=Sat, 0 or 7=Sun). This is what gives the
         wake-catch-up guarantee: if the Mac is asleep at the anchor time on
         an active day, launchd runs the job as soon as it wakes. slap.py
         runner itself then rolls a random moment within schedule.
         fire_window_start/end and sleeps until it, or fires immediately if
         that moment already passed. Plain cron does NOT catch up on wake —
         do not use it. slap.py runner also re-checks active_days itself at
         drain time (runner.is_active_day) — a second line of defense if this
         plist is ever stale relative to config.yaml. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/slap/runner.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/slap/runner.err.log</string>
</dict>
</plist>
